# supervisor/supervisord.conf
[supervisord]
nodaemon=true

[program:storage]
command=/usr/local/bin/gunicorn -c /code/djangows/gunicorn_conf.py --chdir /code/djangows djangorest.wsgi:application --reload
