# supervisor/supervisord.conf
[supervisord]
nodaemon=true

[program:storage-service]
command=/usr/local/bin/gunicorn -c gunicorn_conf.py --chdir /code/djangorest djangorest.wsgi:application --reload
